# deploy/stock-screener.service
[Unit]
Description=Stock Screener (Flask + gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/stock-screener
Environment="PATH=/opt/stock-screener/venv/bin"
# Secrets live in /opt/stock-screener/.env (chmod 600, owner www-data).
# Loaded both by python-dotenv at app startup and by systemd here.
EnvironmentFile=-/opt/stock-screener/.env
ExecStart=/opt/stock-screener/venv/bin/gunicorn -c gunicorn.conf.py app:app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
